# etc/udev/rules.d/10-bluetooth.rules
SUBSYSTEM=="bluetooth", RUN+="/usr/local/bin/bluetooth-udev"
